# service/squid/run
#!/bin/bash

/usr/sbin/squid -YC -f /etc/squid/squid.conf
